# packaging/linux/deb/prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade|deconfigure)
        if [ -d /run/systemd/system ]; then
            systemctl stop sacd-vfs.service 2>/dev/null || true
            systemctl disable sacd-vfs.service 2>/dev/null || true
        fi

        # Unmount if still mounted
        CONF="/etc/sacd-vfs/sacd-vfs.conf"
        if [ -f "$CONF" ]; then
            . "$CONF" 2>/dev/null || true
            if [ -n "$MOUNT_DIR" ] && mountpoint -q "$MOUNT_DIR" 2>/dev/null; then
                fusermount3 -uz "$MOUNT_DIR" 2>/dev/null || true
            fi
        fi
        ;;
esac
